# Aixtiv-Symphony/security/com.aixtiv.rogue-agent-detector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aixtiv.rogue-agent-detector</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/as/asoos/aixtiv-symphony/security/rogue-agent-detector.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/as/asoos/aixtiv-symphony</string>
    
    <key>StandardOutPath</key>
    <string>/Users/as/asoos/aixtiv-symphony/security/detector.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/as/asoos/aixtiv-symphony/security/detector-error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>RestartPreventExitStatus</key>
    <integer>0</integer>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>CLOUD_ML_REGION</key>
        <string>us-west1</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>LowPriorityBackgroundIO</key>
    <false/>
    
    <key>Nice</key>
    <integer>-10</integer>
</dict>
</plist>